# config/service-templates/opai-farmos-sync.timer
[Unit]
Description=OPAI farmOS Weekly Sync — backup, update, deploy on Sundays at 4 AM

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
