# config/service-templates/opai-git-sync.timer
[Unit]
Description=OPAI Daily Git Sync — push changes at 11pm

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
