# config/service-templates/opai-wordpress.service
[Unit]
Description=OP WordPress — Multi-site WordPress management
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-wordpress
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-wordpress
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8096
Restart=always
RestartSec=5
TimeoutStopSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/workspace/synced/opai/tools/opai-wordpress/.env

EnvironmentFile=-%t/opai-vault/opai-wordpress.env
[Install]
WantedBy=default.target
