# config/service-templates/opai-helm.service
[Unit]
Description=OPAI HELM — Autonomous Business Runner (port 8102)
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-helm
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-helm
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8102
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/dallas/.nvm/versions/node/v20.19.5/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/workspace/synced/opai/tools/opai-helm/.env

EnvironmentFile=-%t/opai-vault/opai-helm.env
[Install]
WantedBy=default.target
