# config/service-templates/opai-monitor.service
[Unit]
Description=OPAI Monitor — Web Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-monitor
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-monitor
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/workspace/synced/opai/tools/opai-monitor/.env

EnvironmentFile=-%t/opai-vault/opai-monitor.env
[Install]
WantedBy=default.target
