# config/service-templates/opai-docker-cleanup.timer
[Unit]
Description=OPAI Docker Cleanup — daily at 3am

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
