# config/service-templates/opai-chat.service
[Unit]
Description=OPAI Chat - AI Assistant Interface
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-chat
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-chat
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8888
Restart=always
RestartSec=10
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=-/workspace/synced/opai/tools/opai-chat/.env

EnvironmentFile=-%t/opai-vault/opai-chat.env
[Install]
WantedBy=default.target
