# config/service-templates/opai-discord-bot.service
[Unit]
Description=OPAI Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/discord-bridge
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh discord-bridge
ExecStart=/workspace/synced/opai/tools/discord-bridge/start-bot.sh
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

EnvironmentFile=-%t/opai-vault/discord-bridge.env
[Install]
WantedBy=default.target
