# config/service-templates/opai-tasks.service
[Unit]
Description=OPAI Task Control Panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-tasks
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-tasks
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8081
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/workspace/synced/opai/tools/opai-tasks/.env

EnvironmentFile=-%t/opai-vault/opai-tasks.env
[Install]
WantedBy=default.target
