# config/service-templates/opai-portal.service
[Unit]
Description=OPAI Portal — Auth + Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-portal
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-portal
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8090
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/workspace/synced/opai/tools/opai-portal/.env

EnvironmentFile=-%t/opai-vault/opai-portal.env
[Install]
WantedBy=default.target
