# config/service-templates/opai-bx4.service
[Unit]
Description=OPAI Bx4 — BoutaByte Business Bot (AI business intelligence, port 8100)
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/synced/opai/tools/opai-bx4
ExecStartPre=-/workspace/synced/opai/tools/opai-vault/scripts/vault-env.sh opai-bx4
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8100
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/dallas/.nvm/versions/node/v20.19.5/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/workspace/synced/opai/tools/opai-bx4/.env

EnvironmentFile=-%t/opai-vault/opai-bx4.env
[Install]
WantedBy=default.target
